# data/export/systemd/process_master.target.erb
[Unit]
StopWhenUnneeded=true
Wants=<%= process_names.join(' ') %>
PartOf=<%= "#{app}.target" %>
